# blog
#!/bin/bash
#
# blog      Startup script for the blog application
#
# chkconfig: 345 89 11
# description: Kubosch Blog application

app=blog
port=3002

pidfile=/u/apps/$app/shared/log/mongrel.pid

# See how we were called.
case "$1" in
  start)
		if [ -f $pidfile ] && ! ps h `cat $pidfile` > /dev/null
		then
			rm -f $pidfile
		fi
	    mongrel_rails start -p $port -e production -c /u/apps/$app/current -d 1>/u/apps/$app/shared/log/stdout.log 2>/u/apps/$app/shared/log/stderr.log
        ;;
  stop)
        mongrel_rails stop -c /u/apps/$app/current
        ;;
  restart)
        mongrel_rails restart -c /u/apps/$app/current
        ;;
  status)
        ps -ef | grep $app
        ;;
  *)
        echo $"Usage: $prog {start|stop|restart|status}"
        exit 1
esac

exit 0
